# supervisor.conf
[program:nginx]
command = service nginx start -g "daemon off;"
autostart = True

[program:prod_webserver]
command = /bin/bash /app/gunicorn-start.sh                            ; Command to start app
user = root                                                           ; User to run as
stdout_logfile = /app/logs/gunicorn_supervisor.log                    ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding

[program:dev_webserver]
command = python manage.py runserver 0.0.0.0:82
autostart = True
